# secretario/secretario-email.service
[Unit]
Description=El Secretario — Monitor de Email Pinnacle Holdings
Documentation=https://github.com/pinnaclegroupwi/alex-bot
After=network.target alex-bot.service
Wants=network-online.target

[Service]
Type=simple
User=alexuser
WorkingDirectory=/opt/alex-bot
Environment=PATH=/opt/alex-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/alex-bot/venv/bin/python3 /opt/alex-bot/secretario/email_monitor.py --daemon
Restart=always
RestartSec=30
StandardOutput=null
StandardError=journal

# Reiniciar si falla
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
